# openstack_lease_it/config/instance-spy.cron
#!/bin/sh

export PYTHONPATH=/opt/openstack-lease-it/openstack_lease_it
$PYTHONPATH/bin/instance-spy
